# upstart/pyh-alexa-proxy.conf
description "Program Your Home - Amazon Alexa Proxy"
author "Erik Mulder (ewjmulder@yahoo.com)"

start on runlevel [2345]
stop on runlevel [!2345]

setgid pyh
setuid pyh

exec bash /home/pyh/alexa-proxy/upstart/boot.sh
